# tw-daily-email.service
[Unit]
Description=Taskwarrior Daily Email Report
After=network-online.target

[Service]
Type=oneshot
ExecStart=%h/.task/scripts/tw-daily-email
StandardOutput=journal
StandardError=journal
